# simple-touchscreen-button.path
[Unit]
Description=Monitor touchscreen input

[Path]
PathChanged=/dev/input/touchscreen

[Install]
WantedBy=multi-user.target
